# cmd/tgs-server/tgs-server.service
[Unit]
Description=Telegram bot server
After=network.target

[Service]
ExecStart=tgs-server

[Install]
WantedBy=default.target
