# cmd/tgs/tgs.service
[Unit]
Description=Telegram bot server
After=network.target

[Service]
ExecStart=tgs

[Install]
WantedBy=default.target
